# debian/snipcart-webhook-laposte.service
[Unit]
Description=golang render client service

[Service]
ExecStart=/usr/bin/snipcart-webhook-laposte

[Install]
WantedBy=multi-user.target
